# minio.service
[Unit]
Description=Virtual Distributed Ethernet

[Service]
ExecStart=/root/start_minio.sh

[Install]
WantedBy=multi-user.target
